# ops/systemd/devagent.service
[Unit]
Description=DevAgent Core Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# TEMPORARY: running as the operator user for claude CLI access (Option A).
# TODO: switch to a dedicated service user once Option C is implemented (see TODO.md #14).
# Replace YOUR_USERNAME with the actual OS user running this service.
User=YOUR_USERNAME
WorkingDirectory=/srv/devagent
EnvironmentFile=/srv/devagent/.env
Environment=PYTHONUNBUFFERED=1
# Adjust PATH to include the user's local bin and node (nvm) directories.
Environment=PATH=/home/YOUR_USERNAME/.local/bin:/home/YOUR_USERNAME/.nvm/versions/node/current/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env bash -lc 'python3 -m core.matrix_worker'
Restart=on-failure
RestartSec=3
StandardOutput=append:/var/log/devagent/core.log
StandardError=append:/var/log/devagent/core.log

# ── Systemd security hardening ────────────────────────────────────────────────
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
# Adjust ReadWritePaths to include the operator's development directory.
ReadWritePaths=/srv/devagent /var/log/devagent /home/YOUR_USERNAME/development /srv/repos /srv/agent-artifacts /tmp
ProtectHome=read-only
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictRealtime=yes

[Install]
WantedBy=multi-user.target
